# supervisord.conf
[supervisord]
nodaemon=false
logfile=/home/user/webapp/logs/supervisord.log
pidfile=/home/user/webapp/supervisord.pid
childlogdir=/home/user/webapp/logs

[unix_http_server]
file=/home/user/webapp/supervisor.sock

[supervisorctl]
serverurl=unix:///home/user/webapp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:ml-validator]
command=python3 /home/user/webapp/ml-validator-server-v3.py
directory=/home/user/webapp
autostart=true
autorestart=true
stdout_logfile=/home/user/webapp/logs/ml-validator.log
stderr_logfile=/home/user/webapp/logs/ml-validator-error.log
environment=PYTHONUNBUFFERED="1"

[program:webhook]
command=python3 /home/user/webapp/supabase_webhook_v3.5.py
directory=/home/user/webapp
autostart=true
autorestart=true
stdout_logfile=/home/user/webapp/logs/webhook.log
stderr_logfile=/home/user/webapp/logs/webhook-error.log
environment=PYTHONUNBUFFERED="1"
